# srcs/grafana/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:grafana]
autostart=true
autorestart=true
command=/grafana/bin/grafana-server

[program:telegraf]
autostart=true
autorestart=true
command=telegraf
