# srcs/telegraf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
autostart=true
autorestart=true
command=telegraf