# srcs/metallb/ftps/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:ftps]
autostart=true
autorestart=true
command=vsftpd
